# load-testing/chp-node-tf/scripts/shutdown.service
[Unit]
Description=CHP Shutdown Node Service

[Service]
Type=simple
ExecStop=/etc/init.d/shutdown.sh

[Install]
WantedBy=multi-user.target